# crontab
0 7 * * * /app/run_fetch.sh >> /app/data/cron.log 2>&1
